# systemd/chainweb_snapshoot.timer
[Unit]
Description=Daily Snapshot timer

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true
UTC=true

[Install]
WantedBy=timers.target
